# clipboard-bridge.service
[Unit]
Description=Clipboard Bridge (X11 <-> Wayland)
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/clipboard-sync.py
Restart=always
RestartSec=2

[Install]
WantedBy=graphical-session.target
